# config/messenger-worker.ini
[program:messenger-consume]
command=php /Users/sk/Code/tutorials/symfony-messenger/app/bin/console messenger:consume async_priority-high async --time-limit=3600
user=sk
numprocs=2
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d

[supervisorctl]